# infra/systemd/hermes-factory-pool.service
# DEPRECATED por auditoría P0 — no usar runner legacy.
# Este service fue reemplazado por Hermes Gateway externo.
# No habilitar. No ejecutar. No usar scripts/hermes_factory_runner.py.
# Ver docs/factory/RUNTIME_VIGENTE.md.

[Unit]
Description=Hermes Factory Pool (DEPRECATED — DO NOT ENABLE)
ConditionPathExists=/nonexistent/smartpyme-legacy-hermes-factory-pool-disabled

[Service]
Type=oneshot
ExecStart=/bin/false

[Install]
WantedBy=multi-user.target
